# ansible/roles/monitoring/templates/vps-monitor.service.j2
[Unit]
Description=VPS health monitor (single check run)
After=docker.service network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/vps-monitor.py --once
StandardOutput=journal
StandardError=journal
User=root
